# brain/com.nflbrain.ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  NFL-Brain daily ingest — fires run_brain.sh at 7:30am local, unattended.
  Install:
    cp brain/com.nflbrain.ingest.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.nflbrain.ingest.plist
  Disable:
    launchctl unload ~/Library/LaunchAgents/com.nflbrain.ingest.plist
  If your repo isn't at ~/Downloads/Fantasy, edit the run_brain.sh path below.
  (If the Mac is asleep at 7:30, launchd runs the job once shortly after it next wakes.)
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nflbrain.ingest</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/rsbathla/Downloads/Fantasy/brain/run_brain.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>7</integer>
    <key>Minute</key><integer>30</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/rsbathla/Downloads/NFL-Brain/_status/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/rsbathla/Downloads/NFL-Brain/_status/launchd.err.log</string>
</dict>
</plist>
